# third-party/alpine/board/rg35xxsp/overlay/etc/inittab
# Alpine/OpenRC boot sequence for the RG35xxSP image.
::sysinit:/sbin/openrc sysinit --quiet
::wait:/sbin/openrc boot --quiet
::once:/sbin/openrc default --quiet
::ctrlaltdel:/sbin/reboot
::shutdown:/sbin/openrc shutdown
