# third-party/alpine/board/rg35xxsp/overlay/etc/init.d/sd1-expand
#!/sbin/openrc-run

description="Grow the primary SD card partition on first boot"

DISK_DEV="/dev/mmcblk0"
PART_NUM="3"
PART_DEV="/dev/mmcblk0p3"
SDCARD_MNT="/mnt/sdcard"
STATE_DIR="${SDCARD_MNT}/system/minui"
MARKER_DONE="${STATE_DIR}/.expanded"
MARKER_DONE_LEGACY="${SDCARD_MNT}/.expanded"
LOG="/tmp/sdcard-expand.log"
RESERVE_SECTORS="8192"
FDISK_BIN="/usr/sbin/fdisk"
SPLASH_BIN="/usr/bin/bootsplash"
SPLASH_FB="/usr/share/minui/splash.fb"
SPLASH_PID=""

[ -x "$FDISK_BIN" ] || FDISK_BIN="fdisk"

depend() {
	need localmount
	before sdcard-seed-sp
}

log() {
	echo "$(date +"%F %T") $*" >> "$LOG"
}

ensure_sdcard_mounted() {
	mkdir -p "$SDCARD_MNT"
	mountpoint -q "$SDCARD_MNT" >/dev/null 2>&1 && return 0
	mount "$SDCARD_MNT" >>"$LOG" 2>&1 || true
	mountpoint -q "$SDCARD_MNT" >/dev/null 2>&1
}

has_expand_marker() {
	[ -f "$MARKER_DONE" ] && return 0
	[ -f "$MARKER_DONE_LEGACY" ] && return 0
	return 1
}

show_firstboot_expand_notice() {
	if [ -x "$SPLASH_BIN" ]; then
		"$SPLASH_BIN" --fb /dev/fb0 --image "$SPLASH_FB" --unblank --overlay firstboot-expand --animate >/dev/null 2>&1 &
		SPLASH_PID="$!"
		return 0
	fi
	if [ -r "$SPLASH_FB" ] && [ -w /dev/fb0 ]; then
		[ -w /sys/class/graphics/fb0/blank ] && echo 0 > /sys/class/graphics/fb0/blank 2>/dev/null || true
		dd if="$SPLASH_FB" of=/dev/fb0 bs=4096 2>/dev/null || true
	fi
}

stop_firstboot_expand_notice() {
	if [ -n "$SPLASH_PID" ]; then
		kill "$SPLASH_PID" >/dev/null 2>&1 || true
		wait "$SPLASH_PID" >/dev/null 2>&1 || true
		SPLASH_PID=""
	fi
}

trap stop_firstboot_expand_notice EXIT

wait_for_part() {
	i=0
	while [ "$i" -lt 10 ]; do
		[ -b "$PART_DEV" ] && return 0
		sleep 1
		i=$((i + 1))
	done
	return 1
}

part_start_sector() {
	cat "/sys/class/block/$(basename "$PART_DEV")/start" 2>/dev/null || true
}

part_size_sectors() {
	cat "/sys/class/block/$(basename "$PART_DEV")/size" 2>/dev/null || echo 0
}

partx_disk_geometry() {
	partx -g -o START,END --nr "$PART_NUM" "$DISK_DEV" 2>/dev/null | awk '
		NR == 1 && $1 == "START" { next }
		NF >= 2 { print $1 " " $2; exit }
	'
}

on_disk_start_sector() {
	geom="$(partx_disk_geometry 2>/dev/null || true)"
	if [ -n "$geom" ]; then
		printf '%s\n' "$geom" | awk '{print $1; exit}'
		return 0
	fi
	# fdisk fallback for older partx output/tooling.
	"$FDISK_BIN" -l "$DISK_DEV" 2>/dev/null | awk -v p="$PART_DEV" '
		$1 == p { print ($2 == "*") ? $3 : $2; exit }
	'
}

on_disk_end_sector() {
	geom="$(partx_disk_geometry 2>/dev/null || true)"
	if [ -n "$geom" ]; then
		printf '%s\n' "$geom" | awk '{print $2; exit}'
		return 0
	fi
	"$FDISK_BIN" -l "$DISK_DEV" 2>/dev/null | awk -v p="$PART_DEV" '
		$1 == p { print ($2 == "*") ? $4 : $3; exit }
	'
}

detect_fstype() {
	fstype=""

	if command -v blkid >/dev/null 2>&1; then
		blkid_out="$(blkid "$PART_DEV" 2>/dev/null || true)"
		case "$blkid_out" in
			*'TYPE="'*)
				fstype="$(printf '%s\n' "$blkid_out" | sed -n 's/.*TYPE="\([^"]*\)".*/\1/p' | head -n1)"
				;;
		esac
	fi

	if [ -z "$fstype" ]; then
		fstype="$(awk -v p="$PART_DEV" '$1 == p { print $3; exit }' /proc/mounts 2>/dev/null || true)"
	fi

	if [ -z "$fstype" ]; then
		fstype="$(awk -v m="$SDCARD_MNT" '$2 == m { print $3; exit }' /proc/mounts 2>/dev/null || true)"
	fi

	printf '%s' "$fstype"
}

expand_partition_table() {
	start_sector="$1"
	# Keep the start fixed and accept default last sector to grow to disk end.
	printf 'd\n%s\nn\np\n%s\n%s\n\nw\n' "$PART_NUM" "$PART_NUM" "$start_sector" | "$FDISK_BIN" "$DISK_DEV" >>"$LOG" 2>&1 || true
}

partx_update_partition() {
	if partx -u --nr "$PART_NUM" "$DISK_DEV" >>"$LOG" 2>&1; then
		return 0
	fi

	# Some kernels need delete+add for one partition slot.
	partx -d --nr "$PART_NUM" "$DISK_DEV" >>"$LOG" 2>&1 || true
	partx -a --nr "$PART_NUM" "$DISK_DEV" >>"$LOG" 2>&1
}

start() {
	if ! ensure_sdcard_mounted; then
		log "sdcard mount unavailable; skipping expand"
		return 0
	fi
	if has_expand_marker; then
		return 0
	fi
	log "first-boot sdcard expand start"

	if ! command -v "$FDISK_BIN" >/dev/null 2>&1 || ! command -v partx >/dev/null 2>&1 || ! command -v resize.f2fs >/dev/null 2>&1; then
		log "required tools missing ($FDISK_BIN/partx/resize.f2fs)"
		return 0
	fi

	fstype="$(detect_fstype)"
	if [ -z "$fstype" ]; then
		log "unable to detect fs type for $PART_DEV; assuming f2fs"
		fstype="f2fs"
	fi
	if [ "$fstype" != "f2fs" ]; then
		log "unexpected fs type on $PART_DEV: '$fstype' (expected f2fs)"
		return 0
	fi

	start_sector="$(part_start_sector)"
	cur_size="$(part_size_sectors)"
	disk_sectors="$(cat /sys/class/block/mmcblk0/size 2>/dev/null || true)"
	if [ -z "$start_sector" ] || [ -z "$cur_size" ] || [ -z "$disk_sectors" ]; then
		log "unable to read partition geometry"
		return 0
	fi

	target_end=$((disk_sectors - RESERVE_SECTORS - 1))
	if [ "$target_end" -le "$start_sector" ]; then
		target_end=$((disk_sectors - 1))
	fi
	expected_min=$((target_end - start_sector + 1))
	kernel_size="$(part_size_sectors)"

	show_firstboot_expand_notice
	if [ "$kernel_size" -lt "$expected_min" ]; then
		on_disk_start="$(on_disk_start_sector)"
		on_disk_end="$(on_disk_end_sector)"

		if [ "$on_disk_start" != "$start_sector" ] || [ -z "$on_disk_end" ] || [ "$on_disk_end" -lt "$target_end" ]; then
			log "expanding partition table (start=$start_sector target~=$target_end)"
			expand_partition_table "$start_sector"
			on_disk_start="$(on_disk_start_sector)"
			on_disk_end="$(on_disk_end_sector)"
		fi

		if [ "$on_disk_start" != "$start_sector" ] || [ -z "$on_disk_end" ] || [ "$on_disk_end" -lt "$target_end" ]; then
			log "partition table update did not persist (start=$on_disk_start end=$on_disk_end target~=$target_end)"
			return 0
		fi

		log "updating kernel partition entry via partx"
		if ! partx_update_partition; then
			log "partx update failed"
			return 0
		fi

		mdev -s >/dev/null 2>&1 || true
		if ! wait_for_part; then
			log "partition node did not reappear"
			return 0
		fi

		kernel_size="$(part_size_sectors)"
		if [ "$kernel_size" -lt "$expected_min" ]; then
			log "kernel still uses old partition size after partx (kernel=$kernel_size expected_min=$expected_min)"
			return 0
		fi
	fi

	if mountpoint -q "$SDCARD_MNT"; then
		if ! umount "$SDCARD_MNT" >>"$LOG" 2>&1; then
			log "unable to unmount $SDCARD_MNT - cannot resize"
			return 0
		fi
	fi
	log "resizing f2fs on $PART_DEV"
	if ! resize.f2fs "$PART_DEV" >>"$LOG" 2>&1; then
		log "resize.f2fs failed"
		mount "$SDCARD_MNT" >>"$LOG" 2>&1 || true
		return 0
	fi
	mount "$SDCARD_MNT" >>"$LOG" 2>&1 || true
	mkdir -p "$STATE_DIR"
	touch "$MARKER_DONE"
	log "first-boot sdcard expand complete"
	stop_firstboot_expand_notice
}
