# third-party/alpine/board/rg35xxsp/overlay/etc/init.d/sdcard-seed-sp
#!/sbin/openrc-run
#
# Prepare sdcard-backed SSH and Wi-Fi inputs for the stock Alpine services.
# This keeps first-boot provisioning simple without wrapping the daemons.
#

description="Seed SSH keys and Wi-Fi profiles from the primary SD card"

ssh_dir="/mnt/sdcard/ssh"
ssh_key_path="${ssh_dir}/authorized_keys"
iwd_seed_dir="/mnt/sdcard/wifi/iwd"
iwd_state_dir="/var/lib/iwd"

depend() {
  need localmount
  after localmount sd1-expand
}

seed_authorized_keys() {
  mkdir -p /root/.ssh

  if [ -s "${ssh_key_path}" ]; then
    mkdir -p "${ssh_dir}"
    chmod 600 "${ssh_key_path}" 2>/dev/null || true
    ln -snf "${ssh_key_path}" /root/.ssh/authorized_keys
    return 0
  fi

  rm -f /root/.ssh/authorized_keys
}

seed_iwd_profiles() {
  mkdir -p "${iwd_state_dir}"
  chmod 700 "${iwd_state_dir}" 2>/dev/null || true
  rm -f "${iwd_state_dir}"/*.psk "${iwd_state_dir}"/*.8021x

  [ -d "${iwd_seed_dir}" ] || return 0

  for profile_path in "${iwd_seed_dir}"/*.psk "${iwd_seed_dir}"/*.8021x; do
    [ -f "${profile_path}" ] || continue
    cp -f "${profile_path}" "${iwd_state_dir}/"
  done

  chmod 600 "${iwd_state_dir}"/*.psk "${iwd_state_dir}"/*.8021x \
    2>/dev/null || true
}

start() {
  seed_authorized_keys
  seed_iwd_profiles
}
